# deploy/com.the-rooks-nest.environment-repository.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.the-rooks-nest.environment-repository</string>

    <key>WorkingDirectory</key>
    <string>__APP_ROOT__</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>__APP_ROOT__/deploy/run-production.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__APP_ROOT__/logs/environment-repository.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__APP_ROOT__/logs/environment-repository.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>__HOME__</string>
      <key>PORT</key>
      <string>18082</string>
      <key>ENVIRONMENT_REPOSITORY_LAUNCHD_LABEL</key>
      <string>com.the-rooks-nest.environment-repository</string>
    </dict>
  </dict>
</plist>
